# data/crontab
0 2,14 * * * /root/m0/ctrl-linux/backup.sh >> /tmp/cron.log 2>&1
10 2 * * * /root/m1/ctrl-linux/backup.sh >> /tmp/cron.log 2>&1

